# .container/profile.d/z11_lager_env.sh
#!/bin/bash

## nothing here currently

